# app/defaults/crontab
# Atlas Crontab (supercronic)
# Timezone via TZ env variable in the container

# Daily session cleanup at 06:00 (special: uses --max-turns and session resume)
0 6 * * *  /atlas/app/triggers/cron/daily-cleanup.sh

