# app/defaults/crontab
# Atlas Crontab (supercronic)
# Zeitzone via TZ env-Variable im Container

# Täglicher Session-Cleanup um 06:00 (spezial: nutzt --max-turns und Session-Resume)
0 6 * * *  /atlas/app/triggers/cron/daily-cleanup.sh

